# docker/backend/supervisord.conf
[supervisord]
nodaemon=true
logfile=/app/logs/supervisord.log
pidfile=/var/run/supervisord.pid
user=root

[program:redis]
command=redis-server --port 6379 --bind 127.0.0.1 --daemonize no
autostart=true
autorestart=true
stderr_logfile=/app/logs/redis.err.log
stdout_logfile=/app/logs/redis.out.log
priority=10
user=root

[program:flask]
command=gunicorn --bind 0.0.0.0:5000 --workers 4 --threads 2 --timeout 120 --access-logfile /app/logs/flask.access.log --error-logfile /app/logs/flask.error.log app:create_app()
directory=/app
autostart=true
autorestart=true
stderr_logfile=/app/logs/flask.err.log
stdout_logfile=/app/logs/flask.out.log
environment=FLASK_ENV="production"
user=root

[program:celery]
command=celery -A celery_worker.celery worker --loglevel=info --logfile=/app/logs/celery.log
directory=/app
autostart=true
autorestart=true
stderr_logfile=/app/logs/celery.err.log
stdout_logfile=/app/logs/celery.out.log
user=root

[unix_http_server]
file=/var/run/supervisor.sock

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

